# systemd/hyprfresh.service
[Unit]
Description=HyprFresh - Wayland screensaver daemon for Hyprland
Documentation=https://github.com/solotitan/hyprfresh
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/hyprfresh
Restart=on-failure
RestartSec=5

# Environment
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
